# deploy/launchpad.service
# Launchpad dashboard systemd unit.
#
# Runs the dashboard continuously with the e-ink display driver on the
# Waveshare 7.5" V2 panel (800x480).
#
# Prerequisites:
#   - /opt/launchpad checkout with provisioned venv
#   - Waveshare e-Paper library at /opt/e-Paper/ (symlinked into venv)
#   - aisha user in spi and gpio groups
#   - SPI enabled via raspi-config
#
# Install:
#   sudo cp deploy/launchpad.service /etc/systemd/system/launchpad.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now launchpad
#
# Inspect:
#   systemctl status launchpad
#   journalctl -u launchpad -f

[Unit]
Description=Launchpad e-ink family dashboard
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=aisha
Group=launchpad
WorkingDirectory=/opt/launchpad
Environment=PYTHONPATH=src
Environment=LAUNCHPAD_RUN_FOREVER=1
Environment=LAUNCHPAD_DISPLAY_DRIVER=eink
ExecStart=/opt/launchpad/.venv/bin/python3 -m launchpad
Restart=on-failure
RestartSec=30

# Allow the process to access SPI and GPIO devices
SupplementaryGroups=spi gpio

[Install]
WantedBy=multi-user.target
